# config/s6/sv/bell/up
#!/usr/bin/execlineb -P

fdmove -c 2 1
s6-envdir /home/simen/.local/config/s6/env
multisubstitute
{
	importas DISPLAY DISPLAY
	importas XAUTHORITY XAUTHORITY
}
exec -c
export DISPLAY ${DISPLAY}
export XAUTHORITY ${XAUTHORITY}
xset b off
