# config/s6/sv/sfeed-update-srv/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
s6-envdir /home/simen/.local/config/s6/env
multisubstitute
{
    importas XDG_CONFIG_HOME XDG_CONFIG_HOME
    importas XDG_DATA_HOME XDG_DATA_HOME
}
exec -c
export XDG_CONFIG_HOME ${XDG_CONFIG_HOME}
export XDG_DATA_HOME ${XDG_DATA_HOME}
loopwhilex snooze -v -H* -M /30 -S00 -R180 -s1800
sfeed_update ${XDG_CONFIG_HOME}/sfeed/sfeedrc
