# config/s6/sv/dwm-statusbar/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
s6-envdir /home/simen/.local/config/s6/env
multisubstitute
{
	importas DISPLAY DISPLAY
	importas XAUTHORITY XAUTHORITY
	importas XDG_RUNTIME_DIR XDG_RUNTIME_DIR
	importas XDG_CONFIG_HOME XDG_CONFIG_HOME
    importas DBUS_SESSION_BUS_ADDRESS DBUS_SESSION_BUS_ADDRESS
}
exec -c
export DISPLAY ${DISPLAY}
export XAUTHORITY ${XAUTHORITY}
export XDG_RUNTIME_DIR ${XDG_RUNTIME_DIR}
export XDG_CONFIG_HOME ${XDG_CONFIG_HOME}
export DBUS_SESSION_BUS_ADDRESS ${DBUS_SESSION_BUS_ADDRESS}
${XDG_CONFIG_HOME}/s6/sv/dwm-statusbar/statusbar.sh
